# provision/roles/videotagloader/files/videotag.service
[Unit]
Description=videotag

[Service]
User=videotag
Group=videotag
WorkingDirectory=/opt/ipfs-video-tags
ExecStart=/opt/ipfs-video-tags/ipfs-video-tags-exe
Restart=always

[Install]
WantedBy=multi-user.target
